# config/crontab
# Run every hour
0 * * * * /app/xtrabackup.sh
